# remaudiocd.service
[Unit]
Description=audio cd ejection

[Service]
Type=simple
ExecStart=/usr/local/bin/remaudiocd.sh

[Install]
WantedBy=multi-user.target
